# contrib/gentoo/nield
#!/sbin/runscript

CMD="/usr/sbin/nield"
#extra_commands="-d /var/log/nield.dbg"

[ -x ${CMD} ] || exit 1

start () {
	ebegin "Starting Network Interface Events Logging Daemon"
	start-stop-daemon --start --exec ${CMD} -- ${extra_commands}
	eend $?
}

stop() {
	ebegin "Stopping Network Interface Events Logging Daemon"
	start-stop-daemon --stop --exec ${CMD}
	eend $?
}
